# deploy/systemd/omnidesk-maintenance.timer
[Unit]
Description=Run OmniDesk maintenance daily

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
