# deploy/systemd/omnidesk-maintenance.service
[Unit]
Description=Run OmniDesk SQLite and disk maintenance

[Service]
Type=oneshot
User=omnidesk
Group=omnidesk
ExecStart=/opt/omnidesk/venv/bin/python /opt/omnidesk/scripts/maintenance_sqlite.py --integrity-check --wal-checkpoint --vacuum /var/lib/omnidesk/*.sqlite3
ExecStartPost=/opt/omnidesk/venv/bin/python /opt/omnidesk/scripts/check_disk_guard.py --path /var/lib/omnidesk --max-used-percent 85 --min-free-mb 1024
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/omnidesk
