# deploy/sandbox-runner.service
[Unit]
Description=OmniDesk Remote Sandbox Runner
Documentation=file:/opt/omnidesk/docs/SANDBOX_RUNNER_PRODUCTION.md
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/omnidesk/sandbox-runner.env
WorkingDirectory=/opt/omnidesk
ExecStart=/opt/omnidesk/.venv/bin/python -m omnidesk_agent.sandbox.runner_server
Restart=on-failure
RestartSec=5
User=omnidesk-sandbox
Group=omnidesk-sandbox
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/omnidesk-sandbox-workspaces /tmp
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
